# Scripts/crontab
1 * * * * run-backups home full bootable  
